# bin/docker-mz8.service
[Unit]
Description=iperf3 Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=docker-all-containers start
Restart=on-abort


[Install]
WantedBy=multi-user.target
