# bin/docker-mz8.service
[Unit]
Description=iperf3 Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/sbin/docker-bring-up
Restart=on-abort


[Install]
WantedBy=multi-user.target
