# bin/docker-mz8.service
[Unit]
Description=Docker MZ8 Service
Requires=docker.service
After=docker.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/docker-all-containers start
Restart=on-abort

[Install]
WantedBy=multi-user.target
